# _site/deploy/systemd/geeknews-oneshot.service
[Unit]
Description=GeekNews Blog Automation (One-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/my-blog-cli
Environment="PATH=/home/ubuntu/my-blog-cli/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/ubuntu/my-blog-cli/.env

# 1회 실행
ExecStart=/home/ubuntu/my-blog-cli/venv/bin/python /home/ubuntu/my-blog-cli/scripts/run_once.py

# 로그 설정
StandardOutput=append:/home/ubuntu/my-blog-cli/logs/oneshot.log
StandardError=append:/home/ubuntu/my-blog-cli/logs/oneshot.error.log

# 리소스 제한
MemoryMax=512M
CPUQuota=50%


